# 99-olkb-planck.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="feed", ATTRS{idProduct}=="6060", RUN+="/home/james/code/pollkb.sh"
